# fusionpbx/supervisord.conf
[supervisord]
nodaemon=true

[program:php5-fpm]
command=/usr/sbin/php5-fpm
priority=200

[program:nginx]
command=/usr/sbin/nginx
priority=400
